# bdprochot_off.service
[Unit]
Description=Turns off bdproc

[Service]
ExecStart=/usr/bin/bdprochot.sh off
Type=oneshot

[Install]
WantedBy=default.target
